# tools/gtlsshd.service.in
[Unit]
Description=Remote login daemon from gtlssh
Requires=network.target
Documentation=man:gtlsshd(8)

[Service]
Type=simple
EnvironmentFile=@sysconfigdir@/default/gtlsshd
ExecStart=@gtlsshdbindir@/gtlsshd --nodaemon $gtlsshdflags
StandardOutput=null
Restart=on-failure

[Install]
WantedBy=multi-user.target
